# lifebook/suspend-login.service
[Unit]
Description=Suspend Login

[Service]
ExecStart=/usr/local/bin/suspend-login.sh
Restart=always

[Install]
WantedBy=multi-user.target

